# 70-persistent-net.rules
SUBSYSTEM=="net", ACTION=="add", ATTR{address}=="40:a5:ef:f2:09:55", NAME="wlan1", RUN+="/sbin/iwconfig wlan1 power off"